# scripts/com.prismis.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prismis.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>DAEMON_PATH_PLACEHOLDER</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StartInterval</key>
    <integer>1800</integer>
    
    <key>StandardOutPath</key>
    <string>LOG_PATH_PLACEHOLDER/prismis-daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>LOG_PATH_PLACEHOLDER/prismis-daemon.error.log</string>
    
    <key>WorkingDirectory</key>
    <string>HOME_PATH_PLACEHOLDER</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>HOME_PATH_PLACEHOLDER</string>
        <key>PATH</key>
        <string>HOME_PATH_PLACEHOLDER/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>